# dev/modules/supervisord/templates/program.conf.erb
;
; /etc/supervisord/conf.d/<%=@name%>.conf 
; This file is generated by Puppet
;
[program:<%= @name %>]
command=<%= @command %>
process_name=<%= @process_name %>
numprocs=<%= @numprocs %>
directory=<%= @directory %>
umask=<%= @umask %>
priority=<%= @priority %>
autostart=<%= @autostart %>
autorestart=<%= @autorestart %>
startsecs=<%= @startsecs %>
startretries=<%= @startretries %>
exitcodes=<%= @exitcodes %>
stopsignal=<%= @stopsignal %>
stopwaitsecs=<%= @stopwaitsecs %>
user=<%= @user %>
redirect_stderr=<%= @redirect_stderr %>
stdout_logfile=<%= @stdout_logfile %>
stdout_logfile_maxbytes=<%= @stdout_logfile_maxbytes %>
stdout_logfile_backups=<%= @stdout_logfile_backups %>
stdout_capture_maxbytes=<%= @stdout_capture_maxbytes %>
stdout_events_enabled=<%= @stdout_events_enabled %>
stderr_logfile=<%= @stderr_logfile %>
stderr_logfile_maxbytes=<%= @stderr_logfile_maxbytes %>
stderr_logfile_backups=<%= @stderr_logfile_backups %>
stderr_capture_maxbytes=<%= @stderr_capture_maxbytes %>
stderr_events_enabled=<%= @stderr_events_enabled %>
environment=<%= @environment %>
serverurl=<%= @serverurl %>
